# iRedMail/samples/sogo/sogo.cron
# SOGo session cleanup and email reminder, should be run every minute.
*   *   *   *   *   PH_SOGO_CMD_TOOL expire-sessions 30 2>/dev/null; PH_SOGO_CMD_EALARMS_NOTIFY

# SOGo vacation messages expiration
# The credentials file should contain the sieve admin credentials (username:passwd)
0   0   *   *   *   PH_SOGO_CMD_TOOL update-autoreply -p PH_SOGO_SIEVE_CREDENTIAL_FILE 2>/dev/null
