# samples/systemd/mariadb.service.d/override.conf
# Reference: https://mariadb.com/kb/en/the-mariadb-library/server-system-variables/#open_files_limit
[Service]
LimitNOFILE=infinity
